# 2-WAN/Quadlet/Fluxer/fluxer-unfurl-shard.container
[Unit]
Description=Fluxer Unfurl Shard
Requires=fluxer-nats.service
After=fluxer-nats.container

[Install]
WantedBy=default.target

[Service]
Restart=on-failure
SystemCallArchitectures=native
MemoryDenyWriteExecute=false

[Container]
Image=ghcr.io/fluxerapp/fluxer-unfurl:latest
ContainerName=fluxer-unfurl-shard
Timezone=Europe/Warsaw
EnvironmentFile=/srv/fluxer/.env
Environment=FLUXER_SVC_MODE=shard
Environment=FLUXER_SVC_SHARD_ID=0
AutoUpdate=registry
NoNewPrivileges=true
DropCapability=all
Network=5Unfurl.network